# data/asusd.service
[Unit]
Description=ASUS Notebook Control

[Service]
ExecStart=/usr/bin/asusd
Restart=on-failure
Type=dbus
BusName=org.asuslinux.Daemon